# crontab
0 2 * * * /Users/egorageev/Desktop/DB_HUESOSY/backup_script.sh
